# fake-webserver/supervisord.conf
[supervisord]
nodaemon=true

[program:fake-apache]
command=python /app/apache-fake-log-gen.py --output LOG --num 0 --prefix /logs/fake

[program:file-beat]
command=/usr/bin/filebeat -e
